# sys/10-local.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/emikkva/.Xauthority", RUN+="/home/emikkva/bin/mons"
